# .crontab
0 10 * * 1 bash /Users/okunoe01/.rsyncDotFilesEveryWeek
0 16 * * * open -a /Applications/Notion.app
