# .crontab
0 10 * * 1 bash /Users/okunoe01/.rsyncDotFilesEveryWeek
0 8 * * * open -a /Applications/Firefox.app https://jira.dev.bbc.co.uk/secure/RapidBoard.jspa?rapidView=8748&projectKey=SUMMA
