# crontabs
0 10 * * 1 bash /Users/okunoe01/.rsyncDotFilesEveryWeek
